# quoteOfDay.service
[Unit]
Description=Quote of the Day

[Service]
ExecStart=/home/t/.config/quoteOfTheDay/quote.sh 

[Install]
WantedBy=default.target

