# quoteOfDay.timer
[Unit]
Description=Run My Script Every Hour

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target

